# docker/service/nd-iot-service-plf-xiaoai.service
[Unit]
Description=nd-iot-service-plf-xiaoai

[Service]
WorkingDirectory=/services/nd-iot-services
ExecStart=/usr/local/bin/npm run nd:platform:xiaoai:start
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
