# docker/service/nd-iot-service-plf-nd.service
[Unit]
Description=nd-iot-service-plf-nd

[Service]
WorkingDirectory=/services/nd-iot-services
ExecStart=/usr/local/bin/npm run nd:platform:nd:start
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
